# profile.d/.zsh_env
# more macOS/Bash-like word jumps
export WORDCHARS=""
