# profile.d/.zsh_antigen
source /usr/local/share/antigen/antigen.zsh

# Load the oh-my-zsh's library.
antigen use oh-my-zsh

# Bundles from the default repo (robbyrussell's oh-my-zsh).
antigen bundle git
antigen bundle thefuck
antigen bundle last-working-dir

# Syntax highlighting bundle.
# antigen bundle zsh-users/zsh-syntax-highlighting

# Bundles agents
antigen bundle ssh-agent
antigen bundle gpg-agent

# Load the theme.
antigen theme robbyrussell

# Tell Antigen that you're done.
antigen apply

# Reload ZSH Config
reload() {
    if [ `command -v antigen` ]; then
        antigen reset
        return 0
    fi

    test -e ~/.zshrc && source ~/.zshrc
}
